# inference/efficiency/dependencies/efficiency-pentathlon/profile.d/path.sh
# Update path to include commands added via pip.
export PATH="$PATH:$HOME/.local/bin"
