# inference/efficiency/dependencies/efficiency-pentathlon/profile.d/prompt.sh
# Colorize shell prompt and add an icon to make it clear that the
# shell is an interactive session.
export PS1="\[\033[00;35m\][beaker]\[\033[00;32m\] \u \[\033[00;36m\]\w \[\033[00m\]$ "
